# templates/default/upstart.conf.erb
description "statsd service"

start on runlevel [2345]
stop on runlevel [016]
setuid <%= @user %>

console log

env PATH=/usr/bin:/usr/local/bin:<%= @install_path %>/bin
chdir <%= @install_path %>

exec statsd <%= @config_path %>/config.js
